# powersync.service
[Unit]
Description=Power Sync - Tesla Powerwall & Amber Electric Integration
Documentation=https://github.com/your-repo/power-sync
After=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/power-sync

# Environment variables
Environment=FLASK_APP=run.py
Environment=FLASK_ENV=production

# Run database migrations before starting
ExecStartPre=/home/pi/power-sync/venv/bin/flask db upgrade

# Start the Flask application with gunicorn
ExecStart=/home/pi/power-sync/venv/bin/gunicorn --bind 0.0.0.0:5000 --workers 2 --timeout 120 "app:create_app()"

# Restart configuration
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=powersync

[Install]
WantedBy=multi-user.target
